# patterns/device-samsung-maguro-configs/lib/systemd/system/bcm4330-bluetooth-firmware.service
[Unit]
Description=Activate Bluetooth

[Service]
Type=simple
ExecStart=/sbin/bcm4330-bluetooth-activate.sh
