# patterns/device-samsung-maguro-configs/etc/profile.d/gstomx.sh
export GST_OMX_CONFIG_DIR="/etc/gstomx"
